# src/build/package-content/lib/systemd/system/moneyspot.service
[Unit]
Description=MoneySpot
After=network.target

[Service]
ExecStart=dotnet /var/moneyspot/MoneySpot6.WebApp.dll
WorkingDirectory=/var/moneyspot
Restart=on-failure
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target