# src/build/package-content/DEBIAN/prerm
#!/bin/sh
set -e

systemctl stop moneyspot.service
systemctl disable moneyspot.service